# roles/pinger/files/pinger.service
[Unit]
Description=Pinger is reporting connecton latency
Documentation=https://github.com/romankuchin/ica0002-2022/blob/main/08-logging/README.md
After=network-online.target

[Service]
User=pinger
ExecStart=/usr/local/bin/pinger

[Install]
WantedBy=multi-user.target
